# linux/systemd/oneshot.service
[Unit]
Description=sample service

[Service]
ExecStart=/bin/ps
Type=oneshot

[Install]
WantedBy=multi-user.target
